# entware-ng/opt/lib/opkg/info/entware-opt.list
/opt/etc/init.d/rc.unslung
/opt/etc/ld.so.conf
/opt/etc/init.d/rc.func
/opt/etc/profile
